# nexus-updater.timer
[Unit]
Description=Run Nexus Updater every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
